# native/linux/etc/udev/rules.d/10-log-all.rules
ACTION=="add", RUN+="/usr/local/bin/log-udev --output /var/log/udev-all.log"
