# native/linux/etc/udev/rules.d/10-log-usb-storage.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", RUN+="/usr/local/bin/log-udev --output /var/log/udev-disks.log"
